# systemd/aur-cleanup.timer
[Unit]
Description=Run aur-builder every day at 3am

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true

[Install]
WantedBy=timers.target